# supervisor/fridge.conf
[program:fridge]
command=/usr/local/bin/fridge.py
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/fridge.err.log
stdout_logfile=/var/log/fridge.out.log
